# deploy/systemd/ops-backend.service.example
[Unit]
Description=Ops Digital Employee Backend
Wants=ops-vllm.service
After=network.target ops-vllm.service

[Service]
Type=simple
User=ops-employee
Group=ops-employee
WorkingDirectory=/opt/ops-employee/software
EnvironmentFile=/etc/ops-employee/backend.env
ExecStart=/opt/miniconda/envs/soft/bin/python -m uvicorn backend.app.main:app --host 0.0.0.0 --port 8010
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
